# target/linux/mvebu/base-files/etc/uci-defaults/02_network
#!/bin/sh
#
# Copyright (C) 2014 OpenWrt.org
#

[ -e /etc/config/network ] && exit 0

touch /etc/config/network

. /lib/functions/uci-defaults.sh
. /lib/mvebu.sh

ucidef_set_interface_loopback

board=$(mvebu_board_name)

case "$board" in
armada-385-db-ap)
	ucidef_set_interfaces_lan_wan "eth1 eth2" "eth0"
	;;
armada-388-clearfog)
	ucidef_set_interfaces_lan_wan "eth1" "eth0"
	;;
armada-395-db-gp)
	ucidef_set_interfaces_lan_wan "nic2" "nic0"
	;;
*)
	ucidef_set_interface_lan "eth0"
	;;
esac

uci commit network

exit 0
